# .docker_202/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=npm run start:dev
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/app.err.log
stdout_logfile=/var/log/app.out.log
